# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
loglevel=info
user=root

; ── Web UI (uvicorn) ───────────────────────────────────────────────────────
[program:uvicorn]
command=uvicorn app:app --host 0.0.0.0 --port 8080
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

; ── Cron daemon (supercronic) ─────────────────────────────────────────────
[program:supercronic]
command=/usr/local/bin/supercronic /app/crontab
directory=/app
autostart=true
autorestart=true
startsecs=3
startretries=5
stdout_logfile=/app/logs/sync.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stderr_logfile=/app/logs/sync.log
stderr_logfile_maxbytes=0
priority=20
